# docker/php/crontab
* * * * * cd /srv/api && /usr/local/bin/php vendor/bin/crunz schedule:run


