# etc/php/crontab
* * * * * cd /srv/app && /usr/local/bin/php bin/crunz schedule:run


